# infra/supervisor/almanac-embed.conf
; FastAPI embed worker. Two processes from day one — the embed-queue dispatcher
; is designed for horizontal scaling. Supervisord restarts on crash; stderr +
; stdout funnel into per-process logs under /var/log/almanac.

[program:almanac-embed]
command=/var/www/almanac/current/workers/embed/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8001 --workers 1
directory=/var/www/almanac/current/workers/embed
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
user=ubuntu
environment=ALMANAC_EMBED_WORKER_SHARED_SECRET="REPLACE_ME",ALMANAC_EMBED_PROVIDER="mock"
stdout_logfile=/var/log/almanac/embed-%(process_num)02d.out.log
stderr_logfile=/var/log/almanac/embed-%(process_num)02d.err.log
stopasgroup=true
killasgroup=true
priority=20
