# infra/supervisor/almanac-horizon.conf
; Laravel Horizon — manages the embed + default queues for Almanac.
; Process group restarts cleanly; supervisor + signal-handling are
; coordinated by Horizon's process model.

[program:almanac-horizon]
process_name=%(program_name)s
command=php /var/www/almanac/current/artisan horizon
autostart=true
autorestart=true
user=ubuntu
redirect_stderr=true
stdout_logfile=/var/log/almanac/horizon.log
stopwaitsecs=120
priority=10
